# geoip/files/geoip-update.cron
# Automatic update of GeoIP database
25 4 * * 2 root /usr/bin/geoipupdate
